# parser bot icloud/competitor_price_monitor/deploy/competitor-price-monitor.service
[Unit]
Description=Competitor Price Monitor Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pricebot
Group=pricebot
WorkingDirectory=/opt/competitor_price_monitor
EnvironmentFile=/opt/competitor_price_monitor/competitor_price_monitor/.env
ExecStart=/opt/competitor_price_monitor/competitor_price_monitor/.venv/bin/python -m competitor_price_monitor.telegram_bot
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
